# project/puppet/environments/production/modules/selenoid/templates/selenoid.service.erb
[Unit]
Description=Selenoid service

[Service]
User=selenoid
ExecStart=/opt/selenoid/selenoid_linux_amd64 -listen :<%= @port %> -conf /opt/selenoid/browsers.json -disable-docker --capture-driver-logs -limit 5 -video-output-dir /opt/selenoid/video
WorkingDirectory=/opt/selenoid
SuccessExitStatus=143
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target